# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

cd "${APP_DIR}" && exec s6-setuidgid hotio "${APP_DIR}/Ombi" --storage "${CONFIG_DIR}/app"
